# supervisord-quexf.conf
[program:quexf]
command=/start-quexf.sh
numproces=1
autostart=true
autorestart=false
